# static-binaries/builders/meta/files/etc/profile.d/99-ublue.sh
export PATH="$PATH:/usr/bluefin/bin"
XDG_DATA_DIRS="/usr/bluefin/share:$XDG_DATA_DIRS"
export XDG_DATA_DIRS